# modules/hblock/dotfiles/profile.d/hblock.sh
# hBlock DNS blocking via HOSTALIASES
export HOSTALIASES="$HOME/.config/hblock/hosts"
