# modules/hyprland/dotfiles/systemd/user/gnome-keyring-secrets.service
[Unit]
Description=GNOME Keyring Secret Service - Hyprland Only
After=hyprland-session.target dbus.service
PartOf=hyprland-session.target
ConditionEnvironment=WAYLAND_DISPLAY
ConditionEnvironment=XDG_CURRENT_DESKTOP=Hyprland

[Service]
Type=simple
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=%t/bus
ExecStart=/usr/bin/gnome-keyring-daemon --foreground --components=secrets
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=hyprland-session.target
