# modules/niri/dotfiles/systemd/user/niri-bootstrap.service
[Unit]
Description=Niri bootstrap (init + daemons)
Wants=pipewire.service wireplumber.service niri-ready.service dms.service
After=graphical-session.target niri-session.target niri-ready.service pipewire.service wireplumber.service dms.service
PartOf=niri-session.target
ConditionEnvironment=WAYLAND_DISPLAY
ConditionEnvironment=NIRI_SOCKET
ConditionEnvironment=XDG_CURRENT_DESKTOP=niri

[Service]
Type=oneshot
RemainAfterExit=yes
TimeoutStartSec=60
# Set environment variables if needed, though niri-bootstrap script handles defaults
ExecStart=/usr/bin/bash -lc "%h/.local/bin/niri-bootstrap"
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=niri-session.target
