# modules/transmission/dotfiles/systemd/user/transmission.timer
[Unit]
Description=Delayed start timer for Transmission

[Timer]
OnBootSec=60s

[Install]
WantedBy=timers.target
