# modules/fusuma/dotfiles/systemd/user/fusuma.service
[Unit]
Description=Fusuma gesture recognizer
After=hyprland-session.target niri-session.target
PartOf=hyprland-session.target niri-session.target

[Service]
Type=simple
Environment=XDG_RUNTIME_DIR=/run/user/%U
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin
PassEnvironment=WAYLAND_DISPLAY NIRI_SOCKET HYPRLAND_INSTANCE_SIGNATURE HYPRLAND_SOCKET SWAYSOCK XDG_CURRENT_DESKTOP XDG_SESSION_TYPE XDG_SESSION_DESKTOP
ExecStart=fusuma
Restart=on-failure
RestartSec=3

[Install]
WantedBy=hyprland-session.target niri-session.target
