# modules/hblock/dotfiles/systemd/hblock-update.service
[Unit]
Description=hBlock per-user HOSTALIASES updater
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/hblock-update || true
