# modules/flatpak/dotfiles/systemd/user/flatpak-managed-install.timer
[Unit]
Description=Delayed start timer for Flatpak installation

[Timer]
OnBootSec=90s

[Install]
WantedBy=timers.target
